# debian/kzorp-utils.cron.d
5 * * * *       /usr/sbin/kzorp-stats-updater
